# configs/quadlets/ai-stack-heartbeat.timer
# ai-stack-heartbeat.timer — AI Stack worker heartbeat timer
#
# Template for worker nodes. Generated/installed by bootstrap.sh.
# Location: ~/.config/systemd/user/ai-stack-heartbeat.timer
#
# Fires every 30 seconds. Randomized delay ≤ 5s to avoid thundering-herd
# when multiple workers reboot simultaneously.

[Unit]
Description=AI Stack worker node heartbeat (every 30s)
After=network.target

[Timer]
OnCalendar=*:*:0/30
AccuracySec=1s
Persistent=true

[Install]
WantedBy=timers.target
